# notification/crontab
0 10 * * * python ./notify.py >> /logfile